# systemd/stratis-fstab-setup@.service.in
[Unit]
Description=Set up Stratis filesystems in /etc/fstab
Documentation=man:stratisd(8)
DefaultDependencies=no
Requires=stratisd-min-postinitrd.service

[Service]
ExecStart=@UNITEXECDIR@/stratis-fstab-setup %i
Type=oneshot
